# services/selenium-chrome/cdp-fingerprint-agent.conf
[program:cdp-fingerprint-agent]
priority=25
command=/usr/bin/python3 -u /opt/bin/cdp-fingerprint-agent.py
autostart=true
autorestart=true
startsecs=3
startretries=3
redirect_stderr=true
stdout_logfile=/tmp/cdp-fingerprint-agent-supervisor.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
